# com.locus.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.locus.api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/takemi/.local/bin/uv</string>
        <string>run</string>
        <string>--no-sync</string>
        <string>--directory</string>
        <string>/Users/takemi/Projects/locus-project/locus</string>
        <string>python</string>
        <string>/Users/takemi/Projects/locus-project/locus/scripts/api_server.py</string>
    </array>

    <!-- ログイン時に自動起動 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- クラッシュ時に自動再起動 -->
    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/takemi/Projects/locus-project/locus/scripts/api.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/takemi/Projects/locus-project/locus/scripts/api.log</string>
</dict>
</plist>
