# com.locus.news-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.locus.news-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/takemi/Projects/locus-project/locus/scripts/.venv/bin/python</string>
        <string>/Users/takemi/Projects/locus-project/locus/scripts/news_bot.py</string>
        <string>--poll-only</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/takemi/Projects/locus-project/locus/scripts</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TELEGRAM_BOT_TOKEN</key>
        <string></string>
        <key>TELEGRAM_CHAT_ID</key>
        <string></string>
    </dict>

    <!-- クラッシュ時に自動再起動 -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/takemi/Projects/locus-project/locus/scripts/news-bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/takemi/Projects/locus-project/locus/scripts/news-bot.log</string>
</dict>
</plist>
