# com.claude-kb.flush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-kb.flush</string>

    <!--
      Note: /Users/takemi/ is a placeholder; install-plists.sh sed-replaces
      it with $HOME at install time. uv is resolved via PATH (set below).
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/takemi/Projects/locus-project/locus &amp;&amp; exec uv run --no-sync python scripts/flush_periodic.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/takemi/.local/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/takemi</string>
    </dict>

    <!-- 5分ごとに実行 -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/claude-kb-flush.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-kb-flush-error.log</string>
</dict>
</plist>
