# com.claude-kb.monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-kb.monthly</string>

    <!--
      Note: /Users/takemi/ is a placeholder; install-plists.sh sed-replaces
      it with $HOME at install time. uv is resolved via PATH (set below).
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/takemi/Projects/locus-project/locus &amp;&amp; exec uv run --no-sync python scripts/monthly.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/takemi/.local/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/takemi</string>
    </dict>

    <!-- 毎月28〜31日の21:00に実行。monthly.pyが月末かどうかを判定して実際には月末のみ生成する。 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Day</key>
            <integer>28</integer>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Day</key>
            <integer>29</integer>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Day</key>
            <integer>30</integer>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Day</key>
            <integer>31</integer>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/claude-kb-monthly.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-kb-monthly-error.log</string>
</dict>
</plist>
