# com.locus.collect-news.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.locus.collect-news</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/takemi/.local/bin/uv</string>
        <string>run</string>
        <string>--no-sync</string>
        <string>--directory</string>
        <string>/Users/takemi/Projects/locus-project/locus</string>
        <string>python</string>
        <string>/Users/takemi/Projects/locus-project/locus/scripts/collect_news.py</string>
    </array>

    <!-- 毎朝 07:00 に実行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/takemi/.local/bin</string>
        <key>SOCIAL_DATA_API_KEY</key>
        <string></string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/locus-collect-news.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/locus-collect-news-error.log</string>
</dict>
</plist>
