# com.claude-kb.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-kb.weekly</string>

    <!--
      Note: /Users/takemi/ is a placeholder; install-plists.sh sed-replaces
      it with $HOME at install time. uv is resolved via PATH (set below).
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/takemi/Projects/locus-project/locus &amp;&amp; exec uv run --no-sync python scripts/weekly.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/takemi/.local/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/takemi</string>
    </dict>

    <!-- 毎週金曜 21:00 に実行 (Weekday: 0=Sun, 1=Mon, ..., 5=Fri, 6=Sat) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/claude-kb-weekly.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-kb-weekly-error.log</string>
</dict>
</plist>
